# addon/addbin/root/etc/profile.d/addbin.sh
#!/bin/sh
PATH=/var/media/ftp/bin:$PATH
export PATH
