# addon/fail2ban/root/etc/profile.d/python.sh
#!/bin/sh
PYTHONPATH=/usr/lib/python2.7/site-packages/lib
export PYTHONPATH

